# gentoo_initscript
#!/sbin/runscript

depend() {
	need net
	use dns
}

start() {
	ebegin "Starting pica-node"
	# Ensure that we run from a readable working dir, and that we do not
	# lock filesystems when being run from such a location.
	cd /
	start-stop-daemon --start --quiet --background --make-pidfile \
		--pidfile /var/run/pica-node.pid --exec /usr/bin/pica-node \
		--stdout /var/log/pica-node/pica-node.log \
		--stderr /var/log/pica-node/pica-node.err \
		--user ${PICANODE_USER:-pica-node}  -- \
		 ${PICANODE_OPTS:--v}
	eend $?
}

stop() {
	ebegin "Stopping pica-node"
	start-stop-daemon --stop --quiet --pidfile /var/run/pica-node.pid
	eend $?
}
